# supervisord.conf
; ─────────────────────────────────────────────────────────────────────────────
; Ultimate Reels AI Hunter — Supervisor config
; Used only if you prefer supervisord over the entrypoint.sh approach.
; entrypoint.sh is the recommended / default launcher.
; ─────────────────────────────────────────────────────────────────────────────

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor/supervisor.sock

; ── Xpra (full desktop + HTML5 client) ────────────────────────────────────────
; Replaces the old xvfb + fluxbox + x11vnc + websockify + noVNC programs.
[program:xpra]
command=xpra start-desktop %(ENV_DISPLAY)s
    --bind-tcp=0.0.0.0:%(ENV_WEB_PORT)s
    --html=on
    --start=fluxbox
    --daemon=no
    --pixel-depth=%(ENV_SCREEN_DEPTH)s
    --dpi=96
    --desktop-scaling=off
    --opengl=no
    --no-mdns
    --no-notifications
    --no-bell
    --no-pulseaudio
    --speaker=off
    --microphone=off
    --webcam=no
    --clipboard=no
    --file-transfer=no
    --printing=no
autostart=true
autorestart=false
priority=10
stdout_logfile=/var/log/reels-hunter/xpra.log
stderr_logfile=/var/log/reels-hunter/xpra.log
stdout_logfile_maxbytes=10MB
startsecs=5
stopasgroup=true
killasgroup=true

; ── Python agent ──────────────────────────────────────────────────────────────
[program:agent]
command=python3 /app/agent.py
environment=DISPLAY="%(ENV_DISPLAY)s"
autostart=true
autorestart=false
priority=20
stdout_logfile=/var/log/reels-hunter/agent.log
stderr_logfile=/var/log/reels-hunter/agent_err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=20MB
; Give Xpra 8 s to be ready before the agent starts
startsecs=8
stopasgroup=true
killasgroup=true
